# ops/august/august-playtest.service
[Unit]
Description=August Adventure autonomous playtest runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/august-playtest.env
ExecStart=/usr/bin/python3 %h/august_adventure/ops/august/runner.py
