# ops/august/august-playtest.service
[Unit]
Description=Autonomous game playtest runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/august-playtest.env
ExecStart=/usr/bin/bash -lc '/usr/bin/python3 "${AUGUST_REPO_DIR}/ops/august/runner.py"'
